# scripts/systemd/opendeploy-api.service
[Unit]
Description=OpenDeploy API
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=opendeploy
Group=opendeploy
WorkingDirectory=/opt/opendeploy/apps/api
EnvironmentFile=/etc/opendeploy/opendeploy.env
ExecStart=/usr/bin/env node src/server.js
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
